# dfis_redis/cookbook/supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=/usr/bin/redis-server /etc/dfis_redis/redis.conf