# dfis_redis/cookbook/supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=/usr/bin/redis-server /opt/redis/redis.conf